# deploy/systemd-user/northops-schedule-work.service
[Unit]
Description=NorthOps schedule:work
After=network.target
ConditionPathExists=/home/deployer/northops/current

[Service]
Type=simple
WorkingDirectory=/home/deployer/northops/current
ExecStart=/usr/bin/php artisan schedule:work
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
